# systemd/zyxel-nightly-reboot.timer
[Unit]
Description=Nightly reboot of the cabin Zyxel LTE7460 (04:00 Oslo)
Documentation=https://github.com/henriklex/knausen-signal

[Timer]
# Explicit TZ so this fires at 04:00 Oslo regardless of the Pi's
# system timezone. `Persistent=true` catches up if the Pi was off at
# the scheduled time (it'll fire on next boot instead of skipping).
OnCalendar=*-*-* 04:00:00 Europe/Oslo
Persistent=true

[Install]
WantedBy=timers.target
